# beeswax_supervisor.conf
[program:beeswax]
command=/opt/beeswax_gopath/src/senior-project-experiment/senior-project-experiment
directory=/opt/beeswax_gopath/src/senior-project-experiment 
autostart=true
autorestart=true
