# mhn/beeswax_supervisor.conf
[program:beeswax]
command=/opt/beeswax_gopath/src/beeswax/beeswax
directory=/opt/beeswax_gopath/src/beeswax
autostart=true
autorestart=true
